# ansible-configuration/exports/apache_exporter.conf.j2
start on startup

script
    /usr/bin/apache_exporter
end script